# .config/autostart/Solaar.desktop
[Desktop Entry]
Name=solaar
Icon=solaar
Exec=solaar --window hide
Terminal=false
Type=Application
X-GNOME-Autostart-enabled=true
